# scripts/launchd/com.ewb.stockingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ewb.stockingest</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/DEV/Elliott-wave/.venv/bin/python3</string>
        <string>/Users/DEV/Elliott-wave/python/scripts/data_ingest.py</string>
        <string>--scope</string>
        <string>stocks</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/DEV/Elliott-wave</string>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/DEV/Elliott-wave/brain-output/data_ingest.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/DEV/Elliott-wave/brain-output/data_ingest.log</string>
</dict>
</plist>
